# misc_scripts/weatherpi.service
[Unit]
Description=WeatherPi

[Service]
Type=simple
ExecStart=/opt/weather/start-weather.sh

[Install]
WantedBy=multi-user.target
